# .claude-plugin/skills/devops-linux-deploy/assets/systemd-node.service
[Unit]
Description=Node App
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/srv/app/repo
EnvironmentFile=/srv/app/shared/.env
ExecStart=/usr/bin/node /srv/app/repo/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
